# infra/systemd/altyn-market-backup.timer
[Unit]
Description=Run Altyn Market staging backup every night

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
